# systemd/proxyware-watchdog.timer
[Unit]
Description=proxyware 워치독 5분 주기

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s

[Install]
WantedBy=timers.target
